# configs/bbb/pvrinit
#!/sbin/runscript
# add to default runlevel to load SGX drivers


depend() {
	need localmount
	after bootmisc modules
}

start() {
	ebegin "Starting PVR"

	eval modprobe -q pvrsrvkm

	if [ -f /usr/local/bin/pvrsrvctl ] ; then
		/usr/local/bin/pvrsrvctl --start --no-module

		eval modprobe -q omaplfb
	fi		
	eend $?
}


